# docker/supervisor/oods.conf
[program:oods]
command=gunicorn oods.wsgi -c gunicorn_config.py
directory=/app/oods/
autostart=true
autorestart=true
user=root
stopasgroup=true
environment=HTTP_OODS_PORT=9080,HTTPS_OODS_PORT=9443
